# deploy/tech-alert-bot.service
[Unit]
Description=Tech Alert Bot Listener
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/tech-alert-agent
ExecStart=/home/ubuntu/tech-alert-agent/venv/bin/python3 bot_listener.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
